# systemd/forking.service
[Unit]
Description=Trying forking process

[Service]
Type=forking
ExecStart=
Restart=

[Install]
WantedBy=multi-user.target
