# systemd/mywatchfile.path
[Path]
PathChanged=/home/uy/sandbox/systemd/conf

[Install]
WantedBy=multi-user.target
